# docker/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/storage/logs/supervisor.log

[program:cron]
command = crond -f

[program:php]
command = php-fpm